# tellsticklive/rootfs/etc/services.d/cloudsync/finish
#!/command/with-contenv bashio
# ==============================================================================
# Handle cloud sync service finish
# ==============================================================================

# Exit code 0 means the service exited normally (e.g., not enabled)
# In this case, exit 125 to prevent restart
if [[ "${1}" -eq 0 ]]; then
    exit 125
fi

# Exit code 256 is expected for s6-svc -Od .
if [[ "${1}" -eq 256 ]]; then
    exit 125
fi

bashio::log.warning "Cloud sync service crashed with exit code ${1}, restarting..."
