# tellsticklive/rootfs/etc/services.d/cloudsync/run
#!/command/with-contenv bashio

# Cloud Sync Service - Polls Telldus Live API and syncs devices to local config

# Only run if cloud sync is enabled
if ! bashio::config.true "enable_cloud_sync"; then
    bashio::log.info "Cloud sync is disabled. Service not starting."
    exec s6-svc -Od .
fi

# Check if required OAuth credentials are configured
if ! bashio::config.has_value "live_public_key"; then
    bashio::log.error "Cloud sync enabled but live_public_key is not configured"
    exec s6-svc -Od .
fi

if ! bashio::config.has_value "live_private_key"; then
    bashio::log.error "Cloud sync enabled but live_private_key is not configured"
    exec s6-svc -Od .
fi

if ! bashio::config.has_value "live_token"; then
    bashio::log.error "Cloud sync enabled but live_token is not configured"
    exec s6-svc -Od .
fi

if ! bashio::config.has_value "live_token_secret"; then
    bashio::log.error "Cloud sync enabled but live_token_secret is not configured"
    exec s6-svc -Od .
fi

# Wait for telldusd to be ready
bashio::log.info "Waiting for telldusd to be ready..."
SOCKET_WAIT=0
while [[ ! -S /tmp/TelldusClient ]] || [[ ! -S /tmp/TelldusEvents ]]; do
    sleep 1
    SOCKET_WAIT=$((SOCKET_WAIT + 1))
    if [[ ${SOCKET_WAIT} -ge 60 ]]; then
        bashio::log.error "Timeout waiting for telldusd sockets"
        exit 1
    fi
done
bashio::log.info "Telldusd is ready"

# Get configuration values
PUBLIC_KEY=$(bashio::config 'live_public_key')
PRIVATE_KEY=$(bashio::config 'live_private_key')
TOKEN=$(bashio::config 'live_token')
TOKEN_SECRET=$(bashio::config 'live_token_secret')
SYNC_INTERVAL=$(bashio::config 'cloud_sync_interval')

# Default interval if not set
if [[ -z "${SYNC_INTERVAL}" ]] || [[ "${SYNC_INTERVAL}" -lt 60 ]]; then
    SYNC_INTERVAL=300
fi

bashio::log.info "Starting Telldus Live cloud sync service..."
bashio::log.info "Sync interval: ${SYNC_INTERVAL} seconds"

exec python3 /usr/local/bin/cloud_sync.py \
    --public-key "${PUBLIC_KEY}" \
    --private-key "${PRIVATE_KEY}" \
    --token "${TOKEN}" \
    --token-secret "${TOKEN_SECRET}" \
    --config "/etc/tellstick.conf" \
    --interval "${SYNC_INTERVAL}"
